# speeduino-to-mqtt.service
[Unit]
Description=Speeduino ECU to MQTT Bridge
Documentation=https://github.com/askrejans/speeduino-to-mqtt
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
User=speeduino
Group=speeduino

# Path to the binary
ExecStart=/usr/bin/speeduino-to-mqtt --config /etc/speeduino-to-mqtt/settings.toml

# Working directory (app writes no state; /tmp is always available)
WorkingDirectory=/tmp

# Restart policy
Restart=always
RestartSec=10

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
LockPersonality=true

# Serial / TTY access – allow any device the service user can reach.
# The 'dialout' and 'tty' groups give access to /dev/ttyACM*, /dev/ttyUSB*,
# /dev/ttyS*, etc. without hard-coding a specific device path.
# DevicePolicy defaults to 'auto', which inherits the parent cgroup's device
# list and is not restricted further here.
SupplementaryGroups=dialout tty

# Capabilities
CapabilityBoundingSet=
AmbientCapabilities=

# System call filtering
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallErrorNumber=EPERM

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=speeduino-to-mqtt

# Environment variables (optional - can also use settings.toml)
# Environment="SPEEDUINO_LOG_LEVEL=info"
# Environment="SPEEDUINO_MQTT_HOST=mqtt.local"

[Install]
WantedBy=multi-user.target
